# platform/deploy/cog.service
[Unit]
Description=Kiosk — sway Wayland compositor running cog (WPE WebKit)
Documentation=https://github.com/chaddugas/thread_control_panel
After=panel-ui.service systemd-user-sessions.service
Requires=panel-ui.service
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# Replace `pi` with your Pi username if different. The user must be in the
# video, input, and render groups (install-pi.sh handles this).
User=pi
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal

# Architecture: sway is a minimal wlroots-based compositor that owns DRM/KMS
# and applies the output transform that rotates the Waveshare's natively-
# portrait panel (720x1560) to landscape. Cog runs as sway's autostart child
# via WPE-FDO, rendering into a Wayland surface that sway rotates; wlroots
# rotates touch coordinates to match automatically.
#
# We picked sway over cage because Debian Bookworm's `cage` package predates
# the upstream `-r` rotation flag, and there's no apt/cmdline path to rotate
# 90° on Pi OS Lite without a compositor that supports output transforms.
# Sway runs in single-app mode here (see sway-kiosk.config); it isn't doing
# any tiling.
#
# Combined memory footprint (sway + cog) is ~150-180 MB resident on the Pi
# Zero 2 W's 512 MB.
#
# WEBKIT_DISABLE_SANDBOX_THIS_IS_DANGEROUS=1: disables the bubblewrap-based
# sandbox WPE uses to isolate its WebProcess. Required because Debian's
# bubblewrap is shipped with file capabilities but not setuid, which modern
# bwrap rejects ("Unexpected capabilities but not setuid, old file caps
# config?"). We're a kiosk loading only our own localhost UI — the sandbox
# protects against malicious sites we don't have, so disabling it is fine.
Environment=WEBKIT_DISABLE_SANDBOX_THIS_IS_DANGEROUS=1

# Reap any orphan sway/cog from a previous session before starting. Sway's
# PAMName=login opens a login session that can detach child processes
# (notably the cog spawned by sway's `exec` directive) from the unit's
# cgroup, so systemctl stop doesn't always reap them. The leftover cog
# holds the Wayland socket; the next cog instance races it and crashes
# with a `g_object_unref` assertion. The `-` prefix means "don't fail
# the unit if there's nothing to kill" (the common case).
ExecStartPre=-/usr/bin/pkill -9 -x cog
ExecStartPre=-/usr/bin/pkill -9 -x sway

ExecStart=/usr/bin/sway -c /home/pi/thread_control_panel/platform/deploy/sway-kiosk.config

Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
