# platform/deploy/cog.service
[Unit]
Description=Kiosk — cage Wayland compositor running cog (WPE WebKit)
Documentation=https://github.com/chaddugas/thread_control_panel
After=panel-ui.service systemd-user-sessions.service
Requires=panel-ui.service
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# Replace `pi` with your Pi username if different. The user must be in the
# video, input, and render groups (install-pi.sh handles this).
User=pi
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal

# Architecture: cage is a minimal wlroots-based Wayland compositor that owns
# the DRM/KMS device. Cog runs as cage's child, rendering via WPE-FDO into a
# Wayland surface. We need a compositor in the loop because the Waveshare
# 6.25" panel is natively portrait (720x1560) and there's no way to rotate
# the framebuffer 90°/270° without one — Pi OS Lite's vc4-kms-v3d ignores
# `video=...rotate=N` for HDMI outputs. Cage's `-r` flag applies a 90° output
# transform per occurrence (one for 90°, two for 180°, three for 270°), and
# wlroots rotates touch coordinates to match automatically.
#
# Combined memory footprint (cage + cog) is ~140-170 MB resident on the Pi
# Zero 2 W's 512 MB. Cage adds ~10-20 MB on top of cog alone.
#
# WEBKIT_DISABLE_SANDBOX_THIS_IS_DANGEROUS=1: disables the bubblewrap-based
# sandbox WPE uses to isolate its WebProcess. Required because Debian's
# bubblewrap is shipped with file capabilities but not setuid, which modern
# bwrap rejects ("Unexpected capabilities but not setuid, old file caps
# config?"). We're a kiosk loading only our own localhost UI — the sandbox
# protects against malicious sites we don't have, so disabling it is fine.
Environment=WEBKIT_DISABLE_SANDBOX_THIS_IS_DANGEROUS=1
ExecStart=/usr/bin/cage -r -- /usr/bin/cog -P fdo http://localhost:8080

Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
