# platform/deploy/panel-bridge.service
[Unit]
Description=Panel bridge (WebSocket ↔ UART)
Documentation=https://github.com/chaddugas/thread_control_panel
After=network.target

[Service]
Type=simple
# install-pi.sh templates this; tracked source uses `pi` as the placeholder.
User=pi
# Python logs land in the journal in real time, not buffered.
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/panel/current/deploy/panel-bridge.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
